# packaging/linux-leader/files/systemd/pioreactor-web.target
[Unit]
Description=Pioreactor Web (lighttpd + huey)
Wants=lighttpd.service huey.service

[Install]
WantedBy=multi-user.target
